# ansible/roles/claude-runner/templates/claude-runner@.service.j2
[Unit]
Description=Claude Code runner for %i
After=network-online.target
Wants=network-online.target

[Service]
User={{ claude_user }}
Environment=PATH={{ claude_user_home }}/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
EnvironmentFile={{ claude_runner_base_dir }}/env/%i
ExecStart={{ claude_runner_base_dir }}/run.sh ${REPO_PATH}
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
