# ansible/roles/claude-runner/templates/claude-runner@.service.j2
[Unit]
Description=Claude Code runner for %i
After=network-online.target
Wants=network-online.target

[Service]
User={{ claude_user }}
EnvironmentFile={{ claude_runner_base_dir }}/env/%i
WorkingDirectory=${REPO_PATH}
ExecStart={{ claude_runner_base_dir }}/run.sh ${REPO_PATH}
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
