# ansible/roles/runner/templates/claude-runner-api.service.j2
[Unit]
Description=Claude Runner API
After=network.target

[Service]
User={{ claude_user }}
Environment=PATH={{ claude_user_home }}/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=CLAUDE_RUNNER_BASE_DIR={{ runner_base_dir }}
Environment=CLAUDE_RUNNER_PORT={{ runner_port }}
ExecStart=/usr/local/bin/claude-runner-api
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
